# clyvanta-webhook.service
[Unit]
Description=Clyvanta Deployment Webhook Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/clyvanta
ExecStart=/usr/bin/node deploy-webhook.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clyvanta-webhook

# Environment variables
Environment=NODE_ENV=production
Environment=PORT=4040

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/www/clyvanta
PrivateTmp=true

[Install]
WantedBy=multi-user.target